# deploy/ops/railpush-pg-backup.timer
[Unit]
Description=Daily RailPush Postgres backup

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=30m
Persistent=true
Unit=railpush-pg-backup.service

[Install]
WantedBy=timers.target

